# scripts/usb/99-antigrafiti-usb.rules
ACTION=="add", SUBSYSTEM=="block", ENV{DEVTYPE}=="partition", RUN+="/usr/bin/env bash /opt/CERBERUS_ASIST/usb-trigger.sh %E{DEVNAME} add"
