# scripts/usb/99-cerberus-asist-usb.rules
ACTION=="add", SUBSYSTEM=="block", ENV{DEVTYPE}=="partition", RUN+="/usr/bin/env bash /opt/cerberus_asist/usb-trigger.sh %E{DEVNAME} add"
